# services/ouyouyoutube_backend.container
[Container]
ContainerName=ouyouyoutube_backend
# Build au préalable : podman build -t ouyouyoutube_backend:latest ./server
Image=localhost/ouyouyoutube_backend:latest

# Copier server/.env vers ce chemin sur l'hôte (GOOGLE_CLIENT_ID/SECRET, etc.)
EnvironmentFile=%h/.env.d/ouyouyoutube/server.env
Environment=REDIS_URL=redis://ouyouyoutube_redis:6379/0
# `include_granted_scopes=true` (auth.py) fait remonter les scopes accordés
# lors d'une connexion précédente en plus de ceux demandés cette fois-ci —
# oauthlib traite ça comme une erreur par défaut.
Environment=OAUTHLIB_RELAX_TOKEN_SCOPE=1

# --proxy-headers : fait confiance au X-Forwarded-Proto envoyé par Traefik en
# amont, pour que oauthlib voie bien un redirect_uri https en prod.
Exec=uvicorn app.main:app --host 0.0.0.0 --port 8000 --proxy-headers --forwarded-allow-ips=*
# Pas de PublishPort : seul le frontend (Caddy) lui parle, via le réseau/pod.

Network=server_gateway

PodmanArgs=--pod pod_utils

[Service]
Restart=always
RestartSec=10

[Unit]
After=network.target ouyouyoutube_redis.service
Requires=ouyouyoutube_redis.service

[Install]
WantedBy=default.target
